# backend/waitlist.service
[Unit]
Description=N8tive Waitlist API Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/var/www/n8tive/backend
Environment="PATH=/var/www/n8tive/backend/venv/bin"
Environment="STRIPE_SECRET_KEY=sk_live_your_secret_key_here"
Environment="STRIPE_WEBHOOK_SECRET=whsec_your_webhook_secret"
Environment="STRIPE_PRO_PRICE_ID=price_1SPlRyBTJt2ybYLKPxnE2bKY"
Environment="STRIPE_BUSINESS_PRICE_ID=price_1SPlSrBTJt2ybYLKzrlEToNK"
Environment="STRIPE_ENTERPRISE_PRICE_ID=price_1SPlTXBTJt2ybYLKlMMyEwdb"
ExecStart=/var/www/n8tive/backend/venv/bin/python3 /var/www/n8tive/backend/waitlist_api.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
